# ansible/roles/app_acb/templates/app_acb.service.j2
[Unit]
Description={{ service_name }}
After={{ service_after }}

[Service]
User={{ service_user }}
ExecStart={{ service_start_script }}

[Install]
WantedBy=multi-user.target
